# templates/default/unicorn.conf.erb
# Dynamically generated by Chef on <%= node["fqdn"] %>
# Local modifications will be overwritten by Chef.
description "unicorn for <%= @app_name %>"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on runlevel [!12345]

env HOME="/tmp"

pre-start script
  mkdir -p "<%= File.dirname(@socket) %>"
  chown -R <%= @rails_user %>:<%= @rails_group %> "<%= File.dirname(@socket) %>"
end script

respawn

exec start-stop-daemon --start --chdir "<%= @app_path %>" --chuid <%= @rails_user %>:<%= @rails_group %> --exec <%= @bundler ? "#{@bundle_command} exec " : '' %><%= @unicorn_bin %> -- -E <%= @rails_env %> -c "<%= @unicorn_config %>"
